# deploy/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:vector-service]
command=uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app/vector-service
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
environment=WHISPER_MODEL="%(ENV_WHISPER_MODEL)s",WHISPER_DEVICE="%(ENV_WHISPER_DEVICE)s",WHISPER_COMPUTE_TYPE="%(ENV_WHISPER_COMPUTE_TYPE)s"
priority=100

[program:backend]
command=node /app/backend/src/index.js
directory=/app/backend
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
environment=PORT="%(ENV_PORT)s",NODE_ENV="%(ENV_NODE_ENV)s",VECTOR_SERVICE_URL="%(ENV_VECTOR_SERVICE_URL)s",GROQ_API_KEY="%(ENV_GROQ_API_KEY)s",MAX_AUDIO_FILE_MB="%(ENV_MAX_AUDIO_FILE_MB)s"
priority=200

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=300
